# packaging/launchd/com.os1.local.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.os1.local.ollama</string>
  <key>ProgramArguments</key>
  <array>
    <string>__OLLAMA_PATH__</string>
    <string>serve</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key>
    <string>__OLLAMA_BIND__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/ollama.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/ollama.err.log</string>
</dict>
</plist>
