# packaging/launchd/com.os1.local.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.os1.local.health</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/os1-local-ops-health.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>OLLAMA_HOST</key>
    <string>__OLLAMA_HTTP_HOST__</string>
    <key>OS1_LOCAL_OPS_LOG_DIR</key>
    <string>__LOG_DIR__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/local-health.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/local-health.err.log</string>
</dict>
</plist>
